# ressources/timetracker.timer
[Unit]
Description=Appelle le service pour maj des données

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
